# private_dot_config/systemd/user/baseten-anthropic-proxy.service
[Unit]
Description=Baseten Anthropic compatibility proxy for Claude Code

[Service]
Type=simple
ExecStart=%h/.local/bin/baseten-anthropic-fix.py
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1
